# internal/mission/templates/launchagent.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fry.{{.MissionID}}</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{.MissionDir}}/runner.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>{{.IntervalSeconds}}</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>WorkingDirectory</key>
    <string>{{.MissionDir}}</string>
    <key>StandardOutPath</key>
    <string>{{.MissionDir}}/logs/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{.MissionDir}}/logs/launchd.stderr.log</string>
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
